# blackbox.service
[Unit]
Description=Black Box OSINT Scheduler
After=network.target

[Service]
Type=simple
# Customize User, WorkingDirectory, and paths for your installation
User=YOUR_USER
WorkingDirectory=/path/to/blackbox
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=uv run blackbox serve
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
